# srv.service
[Unit]
Description=GoRSS - RSS Reader
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/gorss
ExecStart=/home/exedev/gorss/gorss
Environment=GORSS_DB_PATH=/home/exedev/gorss/data/gorss.db
Environment=GORSS_LISTEN=:8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
